# package/deb/neuron/postinst
#!/bin/bash

echo "postinst"
ldconfig
systemctl daemon-reload
systemctl enable neuron
systemctl start neuron