# package/deb/neuron/prerm
#!/bin/bash

echo "prerm"
systemctl disable neuron
systemctl stop neuron